# set_1_2_integrated__code_sample__GroundTruths_list -- with result.json
[
   "T1021.001: Remote Services: Remote Desktop Protocol (Lateral Movement) -- OilRig RDP to WATERFALLS -- RDP to 10.1.0.6 using xfreerdp\n", 
   "T1021.002: Remote Services: SMB/Windows Admin Shares (lateral-movement) -- Copy Ransomware -- FIN6 utilizes cmd to copy ransomware to an internal distribution server",
   "T1570: Lateral Tool Transfer (lateral-movement) -- Lateral Tool Transfer -- Previously uploaded tools are transferred to bankfileserver using pscp.exe",
   "T1012: Query Registry (discovery) -- Query Registry -- The reg utility is executed via cmd to enumerate a specific Registry key associated with local system policies to ensure that the user will not be prompted for credentials when elevating permissions.",
   "T1087.002: Account Discovery: Domain Account (discovery) -- Enumerate AD person objects -- Find all person objects and output the results to a text file.",
   "T1016: System Network Configuration Discovery (discovery) -- System Network Configuration Discovery -- The ipconfig utility is executed via cmd to enumerate local TCP/IP network configuration information.",
   "T1482: Domain Trust Discovery (discovery) -- Enumerate AD Organizational Units -- Enumerate all Organizational Units (OUs) in the domain of the user running the command and output the results to a text file.",
   "T1082: System Information Discovery (discovery) -- System Information Discovery -- The systeminfo utility is executed via cmd to enumerate local operating system configuration.",
   "T1018: Remote System Discovery (discovery) -- Remote System Discovery -- The net utility is executed via cmd to enumerate DCs within the domain",
   "T1135: Network Share Discovery (discovery) -- Network Share Discovery -- Display network shares",
   "T1518: Software Discovery (discovery) -- Detect Software -- Detect software on host",
   "T1069.002: Permission Groups Discovery: Domain Groups (discovery) -- Enumerate AD groups -- List groups and output the results to a text file.",
   "T1033: System Owner/User Discovery (discovery) -- System Owner / User Discovery -- The native echo command is executed via cmd to enumerate local environment variables associated with current user and domain.",
   "T1049: System Network Connections Discovery (discovery) -- 4.C.2 - System Network Connections Discovery (T1049) -- The netstat utility is executed via cmd to enumerate local active network connections.",
   "T1046: Network Service Scanning (discovery) -- Network Service Scanning -- Enumerate remote hosts for running services using tcping",
   "T1087.001: Account Discovery: Local Account (Discovery) -- OilRig Local Account Discovery -- Enumerate user accounts",
   "T1083: File and Directory Discovery (discovery) -- File and Directory Discovery (T1083) -- PowerShell's Get-ChildItem alias 'ls' is used to enumerate files in a remote file share.",
   "T1069.001: Permission Groups Discovery: Local Groups (Discovery) -- OilRig \"administrators\" local group discovery -- View the details and members of the \"administrators\" local group",
   "T1057: Process Discovery (discovery) -- Process Discovery -- API call(s) are executed to enumerate local running processes.",
   "T1518.001: Software Discovery: Security Software Discovery (discovery) -- Detect Anti-Virus -- Detect anti-virus software on host",
   "T1007: System Service Discovery (discovery) -- System Service Discovery -- The sc utility is executed via cmd to enumerate local active services.",
   "T1136.002: Create Account: Domain Account (Persistence) -- Create a new domain user (Adversary) -- Task Snake to create a new domain user using an access token from one of the domain admin processes. \nThe new domain user will be used as a backdoor domain admin account for persistence on the domain.\n",
   "T1505.003: Server Software Component: Web Shell (persistence) -- OilRig Copy webshell to WATERFALLS -- contact.aspx is copied from THEBLOCK to WATERFALL to the Exchange Web Services directory. This ability implements Step 4.A.3.  \n",
   "T1505.002: Server Software Component: Transport Agent (Persistence) -- Install LightNeuron on the Exchange Server (Adversary) -- Task Snake to install LightNeuron remotely using WMI and PowerShell, using the domain admin's token\n",
   "T1548.002: Abuse Elevation Control Mechanism: Bypass User Account Control (privilege-escalation) -- UAC Bypass via Backup Utility -- Modify registry values of sdclt to bypass UAC",
   "T1548.001: Abuse Elevation Control Mechanism: Setuid and Setgid (privilege-escalation) -- Sandworm Set Up SUID for Agent Execution -- Sandworm executes the downloaded agent executable by running SUID binary /bin/backup.\nThis binary executes check.sh, which causes the agent to be executed with root privileges.\nThis ability must be run by an agent running on the attacker-controlled machine (e.g. on the C2 server itself),\nsince it will target patient zero. Requires curl to be installed on the attacker-controlled machine.\n",
   "T1550.002: Use Alternate Authentication Material: Pass the Hash (Lateral Movement) -- Install the 3rd Snake installer on the third target host (Adversary) -- Task Snake to pass-the-hash using the previously discovered NTLM hash to run PsExec and install Snake on the target workstation.\n",
   "T1547.001: Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder (persistence) -- Set Registry Persistence on CFO -- Set Registry Persistence on CFO to run Java-Update.vbs when the CFO user logs in",
   "T1056.001: Input Capture: Keylogging (collection) -- Automated Collection (T1119) - Input Capture (T1417) -- Load custom PowerShell module, and grab keystrokes for 15 seconds.",
   "T1547.004: Boot or Logon Autostart Execution: Winlogon Helper DLL (persistence) -- Wizard Spider Registry Persistence -- Wizard Spider has established persistence using Userinit by adding the Registry key HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon.\n",
   "T1053.005: Scheduled Task/Job: Scheduled Task (persistence) -- Scheduled Task Persistence -- FIN6 utilizes cmd to execute Scheduled Tasks",
   "T1134.001: Access Token Manipulation: Token Impersonation/Theft (privilege-escalation) -- Bypass User Account Control -- A UAC bypass technique is executed to steal the token of an existing high-integrity process and launch a new, high-integrity RAT with limited functionality.",
   "T1547.009: Boot or Logon Autostart Execution: Shortcut Modification (persistence) -- Persistent Service 1 -- Leverage modified Sysinternals",
   "T1546.011: Application Shimming (Persistence) -- Execute Application Shim Persistence -- Reboot the host to start shim persistence.",
   "T1546.003: Event Triggered Execution: Windows Management Instrumentation Event Subscription (persistence) -- WMI Persistence technique -- user triage",
   "T1055: Process Injection (privilege-escalation) -- Process Injection -- The limited functionality high-integrity RAT will inject malicious code into an existing fully functional high-integrity process, resulting in a new elevated, fully functional high-integrity RAT.",
   "T1543.002: Create or Modify System Process: Systemd Service (persistence) -- Sandworm Establish Systemd Persistence -- Sandworm establishes systemd persistence.\n",
   "T1134.002: Access Token Manipulation: Create Process with Token (defensive-evasion) -- UAC Bypass via sdctl -- Invoke UAC bypass sdctl",
   "T1053.003: Scheduled Task/Job: Cron (persistence) -- Sandworm Establish Crontab Persistence -- Sandworm establishes crontab persistence.\n",
   "T1574.001: Hijack Execution Flow - DLL Search Order Hijacking (Hijack Execution Flow - DLL Search Order Hijacking) -- Privilege Escalation -- Perform DLL hijack to escalate privileges.",
   "T1574.011: Hijack Execution Flow: Services Registry Permissions Weakness (Privilege Escalation) -- Modify the ViperVPN service registry key (Adversary) -- Modify the ViperVPN service registry key\n"
]


**************************************************************************************************************
**************************************************************************************************************

Without-Domain-Knowledge  ( without 'Please list all' )
 - TP ( #7 ) : Remote Services (T1021.001),
               Scheduled Task/Job (T1053.005), 
               Registry Run Keys / Startup Folder (T1547.001), 
               Bypass User Account Control (UAC) (T1548.002), 
               Input Capture (T1056), 
               System Information Discovery (T1082), 
               Account Discovery (T1087)
 - TN: N/A
 - FP ( #10 ): Command and Scripting Interpreter (T1059), 
               User Execution (T1204.002), 
               Modify Existing Service (T1569.002), 
               Valid Accounts (T1078), 
               Unsecured Credentials (T1552), 
               Remote File Copy (T1105), 
               Remote Services (T1021), 
               Data from Information Repositories (T1213),
               Exfiltration Over Alternative Protocol (T1048),
               Application Layer Protocol (T1071)
 - FN ( #34 )

 Accuracy, F-1

############################################################################################
With-Domain-Knowledge  ( without 'Please list all' )
 - TP ( #9 ) : T1021.001 - Remote Services: Remote Desktop Protocol (RDP)
               T1021.002 - Remote Services: SMB/Windows Admin Shares
               T1570 - Lateral Tool Transfer
               T1016 - System Network Configuration Discovery
               T1087.002 - Account Discovery: Domain Account
               T1482 - Domain Trust Discovery
               T1018 - Remote System Discovery
               T1082 - System Information Discovery
               T1547 - Boot or Logon Autostart Execution

 - TN: N/A
 - FP ( #3 ):  T1563.002 - Remote Service Session Hijacking: RDP Hijacking
               T1217 - Browser Bookmark Discovery
               T1059.001 - Command and Scripting Interpreter: PowerShell

 - FN ( #32 )

**************************************************************************************************************
**************************************************************************************************************


Without-Domain-Knowledge  ( WITH 'Please list all' )

 - TP ( #6 ) :  Scheduled Task/Job (T1053.005)
                Registry Run Keys/Startup Folder (T1547.001)
                Event Triggered Execution: Windows Management Instrumentation Event Subscription (T1546.003)
                Abuse Elevation Control Mechanism: Bypass User Access Control (T1548.002)
                Remote Services: Remote Desktop Protocol (T1021.001)
                Input Capture: Keylogging (T1056.001)
 - TN: N/A
 - FP ( #9 ): Command and Scripting Interpreter (T1059)
               System Services: Service Execution (T1569.002)
               Credentials from Password Stores (T1555)
               OS Credential Dumping (T1003)
               Obfuscated Files or Information (T1027)
               Signed Binary Proxy Execution (T1218)
               Remote File Copy (T1105)
               Application Layer Protocol (T1071.001)
               Exfiltration Over C2 Channel (T1041)

 - FN ( #35 )

 Accuracy, F-1

############################################################################################

 With-Domain-Knowledge  ( WITH 'Please list all' )
 - TP ( #9 ) :  T1021.001: Remote Services: Remote Desktop Protocol (RDP)
                T1570: Lateral Tool Transfer
                T1016: System Network Configuration Discovery
                T1082: System Information Discovery
                T1087.002: Account Discovery: Domain Account
                T1033: System Owner/User Discovery
                T1136.001: Create Account: Local Account
                T1543.003: Create or Modify System Process: Windows Service
                T1547.001: Boot or Logon Autostart Execution: Registry Run Keys/Startup Folder

 - TN: N/A
 - FP ( #5 ): T1217: Browser Bookmark Discovery
               T1059.001: Command and Scripting Interpreter: PowerShell
               T1560.002: Archive Collected Data: Archive via Utility
               T1569.002: System Services: Service Execution
               T1204.002: User Execution: Malicious File
 - FN ( #NN )

 Accuracy, F-1